# root/etc/systemd/system/network.service.d/network-recovery-cleanup.conf
[Service]
ExecStartPre=-/usr/bin/bash -c "if ip link show brtmp &>/dev/null ; then echo 'Clean up brtmp bridge' ; exec ip link delete brtmp ; fi"
