# root/etc/systemd/system/network.service.d/network-recovery-cleanup.conf
[Service]
ExecStartPre=-/usr/sbin/ip link delete brtmp
