# root/usr/lib/systemd/system/shorewall.service.d/nethserver-firewall-base.conf
[Unit]
ConditionPathExists=!/var/run/.nethserver-fixnetwork

