# root/etc/cron.d/shorewall-update-dst
10 5 * * * root	/usr/libexec/nethserver/shorewall-update-dst
